# sv-rt.conf
[supervisord]
nodaemon=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:mongos]
command=/usr/bin/mongos __MONGO_OPTIONS
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:init_user]
command=./init_user.sh 10
autostart = true
autorestart = false
startsecs = 0

[program:init_shard]
command=./init_shard.sh 30
autostart = true
autorestart = false
startsecs = 0
